# etc/mbs-update.service
[Unit]
Description=Minecraft Bedrock Server Update

[Service]
Type=oneshot
ExecStart=bash mc/server-setup/bin/update.sh
